# deploy/huis-dashboard.service
[Unit]
Description=Huis Dashboard (LAN dashboard + Ollama agent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/dashboard
# Absolute node path (fnm-managed). Update if the Node version changes.
ExecStart=/root/.local/share/fnm/node-versions/v24.16.0/installation/bin/node /opt/dashboard/src/server.js
Restart=on-failure
RestartSec=3
User=root
# Serve HTTP on 8080 and HTTPS on 443 (cert in /opt/dashboard/data/tls)
Environment=NODE_ENV=production
Environment=DASH_AGENT_MAX_TOOL_CALLS=32

[Install]
WantedBy=multi-user.target
